# bootstrap/files/etc/profile.d/check_aws.sh
test -z "$AWS_ACCESS_KEY" -o -z "$AWS_SECRET_KEY" && echo AWS_ACCESS_KEY and AWS_SECRET_KEY must be defined >&2 && exit 1
