# Docker/root/etc/profile.d/prompt.sh
#!/bin/sh

B='\[\033[01;30m\]'
G='\[\033[01;32m\]'
n='\[\033[0m\]'
b='\[\033[01;34m\]'
export PS1="$B(${G}container$B)$n [\u@\h $b\W$n]\$ "
